# centos/7/supervisord.conf
[supervisord]
nodaemon=true
[program:marklogic]
; pidproxy should does not seem to be working with MarkLogic process
;command=/usr/bin/pidproxy /var/run/MarkLogic.pid /etc/rc.d/init.d/MarkLogic start
command=/bin/bash -c "/etc/rc.d/init.d/MarkLogic start && tail -F /var/opt/MarkLogic/Logs/ErrorLog.txt"
redirect_stderr=true
autorestart=true
startsecs=30 